# systemd/kbdassist-numpad-CTechKBN01.service
[Unit]
Description=kbdassist for CNC pendant

[Service]
Type=simple
ExecStart=/root/kbdassist-wired-KBN01.sh
